# build_routereflector/node_filesystem/etc/rc.local
#!/bin/bash
# Record current boot time. Helps with debugging.
date > /tmp/boottime.txt

# Run Confd in onetime mode, to ensure that we have a working config in place
# to allow bird(s) to start.
#./confd -confdir=. -onetime -node ${ETCD_AUTHORITY}