# debian/calico-bird.bird.upstart
# bird - BIRD Internet Routing Daemon (IPv4)

description "BIRD Internet Routing Daemon (IPv4)"
author "Ondřej Surý <ondrej@debian.org>"

start on runlevel [2345]
stop on runlevel [016]

respawn
pre-start script
/usr/lib/bird/prepare-environment
/usr/sbin/bird -p
end script

script
. /etc/bird/envvars
exec /usr/sbin/bird -f -u $BIRD_RUN_USER -g $BIRD_RUN_GROUP $BIRD_ARGS
end script
